# supervisord.conf
[program:web]
command=/install/gunicorn

[program:nginx]
command=nginx -c /install/nginx.conf

[supervisord]
