# supervisord.conf
[program:web]
command=/install/gunicorn

[supervisord]
